# packages/sysutils/inject_bl301/system.d/update-bl301.service
[Unit]
Description=masQelec BL301 Update Service

[Service]
Type=oneshot
ExecStart=/usr/lib/masqelec/update-bl301 -v
RemainAfterExit=no

[Install]
WantedBy=basic.target
